# mail/deploy/systemd/kgm-mail-service.service
[Unit]
Description=Karacabey Gross Market Mail Service Phase 2
After=network.target docker.service

[Service]
Type=simple
WorkingDirectory=/opt/kgm-mail-service
EnvironmentFile=/opt/kgm-mail-service/.env
ExecStart=/opt/kgm-mail-service/kgm-mail-service
Restart=always
RestartSec=3
User=kgm-mail
Group=kgm-mail
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/opt/kgm-mail-service/data

[Install]
WantedBy=multi-user.target
